# ngrok-monitor.service
[Unit]
Description=Ngrok Monitor + Telegram WiFi Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kali
WorkingDirectory=/home/kali/bot2ssh
ExecStart=/home/kali/bot2ssh/run_monitor.sh
Restart=on-failure
RestartSec=10
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin

[Install]
WantedBy=multi-user.target
